# find-cmdlet-web/pkg/find-cmdlet-web.socket
[Unit]
Description=find-cmdlet-web socket

[Socket]
ListenStream=443
Accept=no

[Install]
WantedBy=sockets.target
